# services/convitelist-update.service
[Unit]
Description=Convitelist Auto-Update Service
After=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/convitelist
# Pull latest code
ExecStart=/usr/bin/git -C /home/pi/convitelist pull origin master
# Rebuild and restart if there were changes
ExecStart=/bin/bash -c 'cd /home/pi/convitelist && if [ -n "$(git diff HEAD@{1} HEAD --name-only 2>/dev/null || echo changed)" ]; then if command -v docker &> /dev/null; then docker compose build && docker compose up -d || docker-compose build && docker-compose up -d; fi; fi'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

