# services/minecraft-sleeper-improved.service
[Unit]
Description=Minecraft Server Sleeper Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/minecraft
Environment="PATH=/home/pi/.nvm/versions/node/v22.21.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Flush memory cache before starting
ExecStartPre=/bin/bash -c 'sudo sh -c "echo 1 > /proc/sys/vm/drop_caches"'
ExecStartPre=/bin/sync
# Fix permissions
ExecStartPre=/bin/bash /home/pi/minecraft/fixpermissions.sh -a
# Start the sleeper interface
ExecStart=/home/pi/.nvm/versions/node/v22.21.1/bin/npm run start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target



